# docker/root/etc/services.d/10-redis/run
#!/usr/bin/execlineb -P

exec redis-server --dir /app/data --bind 0.0.0.0 --logfile /app/logs/redis-server.log
