# docker/root/etc/services.d/10-redis/run
#!/usr/bin/execlineb -P

exec redis-server /etc/redis/redis.conf
