# docker/root/etc/s6-overlay/s6-rc.d/20-hpr/run
#!/usr/bin/with-contenv sh

exec /app/bin/puma --dir /app -C /app/config/puma.rb | tee /app/logs/hpr.log
