# docker/root/etc/services.d/redis/run
#!/usr/bin/with-contenv sh

mkdir -p /data

echo "** Starting Redis.."
exec s6-setuidgid redis redis-server --dir /data
