# docker/root/etc/s6-overlay/s6-rc.d/10-redis/run
#!/command/execlineb -P

redis-server /etc/redis.conf
