# docker/root/etc/services.d/20-hpr/run
#!/usr/bin/with-contenv sh

exec hpr server --path /app 2>&1 | tee -a /app/logs/hpr.log
