# docker/root/etc/services.d/20-hpr/run
#!/usr/bin/with-contenv sh

cd /app
exec hpr server 2>&1 | tee -a /app/logs/hpr.log
