# docker/root/etc/services.d/20-hpr/run
#!/usr/bin/execlineb -P

exec /bin/hpr server --path /app
