# docker/root/etc/s6-overlay/s6-rc.d/21-sidekiq/run
#!/usr/bin/with-contenv sh

exec /app/bin/sidekiq -r /app/lib/hpr.rb -c `cat /proc/cpuinfo | grep -c "cpu cores"`
